# contrib/init/bellsd.conf
description "LuckyCoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env luckycoind_BIN="/usr/bin/luckycoind"
env luckycoind_USER="luckycoin"
env luckycoind_GROUP="luckycoin"
env luckycoind_PIDDIR="/var/run/luckycoind"
# upstart can't handle variables constructed with other variables
env luckycoind_PIDFILE="/var/run/luckycoind/luckycoind.pid"
env luckycoind_CONFIGFILE="/etc/luckycoin/luckycoin.conf"
env luckycoind_DATADIR="/var/lib/luckycoind"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # luckycoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$luckycoind_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run luckycoind."
        echo "The setting must appear in $luckycoind_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"LuckyCoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$luckycoind_PIDDIR"
    chmod 0755 "$luckycoind_PIDDIR"
    chown $luckycoind_USER:$luckycoind_GROUP "$luckycoind_PIDDIR"
    chown $luckycoind_USER:$luckycoind_GROUP "$luckycoind_CONFIGFILE"
    chmod 0660 "$luckycoind_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$luckycoind_PIDFILE" \
    --chuid $luckycoind_USER:$luckycoind_GROUP \
    --exec "$luckycoind_BIN" \
    -- \
    -pid="$luckycoind_PIDFILE" \
    -conf="$luckycoind_CONFIGFILE" \
    -datadir="$luckycoind_DATADIR" \
    -disablewallet \
    -daemon

